# crontab
#run script every friday morning 0:0:0

0 0 * * 5 /usr/local/bin/python /app/app.py > /prov/1/fd/1 2> /prov/1/fd/2
